# supervisord.conf
[supervisord]

[program:renderer]
directory=/var/www/renderer
command=gulp
